# host-setup/systemd/nanoclaw-checkpoint.service
[Unit]
Description=NanoClaw daily checkpoint backup

[Service]
Type=oneshot
ExecStart=/home/pseudo/nanoclaw/scripts/checkpoint.sh
WorkingDirectory=/home/pseudo/nanoclaw
StandardOutput=append:/home/pseudo/nanoclaw/logs/checkpoint.log
StandardError=append:/home/pseudo/nanoclaw/logs/checkpoint.log
